# upnp/etc/network/if-up.d/phpdlna
#!/bin/sh

# start announce with /etc/init.d/phpdlna
/etc/init.d/phpdlna restart
